# deploy/jetson/seestar-stack.service
[Unit]
Description=Seestar S50 Stack (portal + ALP backend)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=jmm2020
WorkingDirectory=/home/jmm2020/seestar-scope
ExecStart=/usr/bin/docker compose --profile tunnel up
ExecStop=/usr/bin/docker compose --profile tunnel down
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=seestar-stack

[Install]
WantedBy=multi-user.target
